# tools/supervisor/conf.d/myblog_celery_worker_supervisor.conf
; ==================================
;  myblog celery worker supervisor
; ==================================

[program:myblog_worker]
; Set full path to celery program if using virtualenv
;command=/home/nico/.envs/myblog_env/bin/celery worker -A MyBlog -l INFO
command=celery worker -A MyBlog -l INFO
; project root directory
;directory=/home/nico/gitee/MyBlog
directory=/MyBlog
;user=nico
numprocs=1
stdout_logfile=/var/log/myblog_worker.log
stderr_logfile=/var/log/myblog_worker_error.log
autostart=true
autorestart=true
startsecs=5

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 60

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; Higher priorities indicate programs that start last and shut down first.
priority=1001
