# tools/myblog_uwsgi_supervisor.conf
; ==================================
;      myblog uwsgi supervisor
; ==================================

[program:myblog_uwsgi]
command=/usr/local/bin/uwsgi --emperor /etc/uwsgi/vassals
; project root directory
directory=/home/nico/gitee/MyBlog
user=nico
; number of process
numprocs=1
stdout_logfile=/var/log/myblog_uwsgi_access.log
stderr_logfile=/var/log/myblog_uwsgi_error.log
autostart=true
autorestart=true
startsecs=5

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 5

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; Higher priorities indicate programs that start last and shut down first.
priority=999
